# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:redis-server]
priority=1
command=redis-server

[program:worker]
priority=2
directory=/repo
command=rq worker
process_name=%(program_name)s-%(process_num)s
numprocs=4

[program:flask]
priority=3
directory=/repo
command=python run.py